# dungeon-api.service
[Unit]
Description=DungeonKeeper Image API
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/ModMail/render-service
ExecStart=/usr/bin/python3 -m uvicorn app:app --host 0.0.0.0 --port 10000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
